# src/backend/configs/supervisor/backend.conf
[program:app-uwsgi]
user = uwsgi
command = /usr/local/bin/uwsgi --ini /tmp/uwsgi/uwsgi.ini
autostart=true
autorestart=true
killasgroup=true
stopsignal=QUIT