# services/backend/configs/supervisor/api.conf
[supervisord]
nodaemon=true

[program:app]
priority=1
user = uwsgi
command = /usr/local/bin/uwsgi --ini /tmp/uwsgi/uwsgi.ini
autorestart=false
